# nym-vpn-core/crates/nym-vpnd/.pkg/aur/nym-vpnd.service
[Unit]
Description=NymVPN daemon
StartLimitBurst=6
StartLimitIntervalSec=24
Wants=network-pre.target
After=network-pre.target NetworkManager.service systemd-resolved.service

[Service]
ExecStart=/usr/bin/nym-vpnd
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

